# container/solr/conf/supervisor/solr.conf
[program:solr]
command=/usr/local/solr7/solr-7.2.1/bin/solr start -force -p 8983 -e "techproducts"
