# container/solr/conf/supervisor/solr.conf
[program:solr]
command=/usr/local/solr7/solr-7.1.0/bin/solr start -force -p 8983 -s "example/testdrive/solr"
